# vagrant/startups/ptp-slave-1.service
[Unit]
Description = PTP Slave 1

[Service]
ExecStart = /opt/ptp-slave-1.sh
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
