# vagrant/startups/ptp-master.service
[Unit]
Description = PTP Master

[Service]
ExecStart = /opt/ptp-master.sh
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
